# scm/mousepad.scm.list
/apps/mousepad/bin/mousepad
/apps/mousepad/local/bin/dbus-binding-tool
/apps/mousepad/local/bin/dbus-cleanup-sockets
/apps/mousepad/local/bin/dbus-daemon
/apps/mousepad/local/bin/dbus-launch
/apps/mousepad/local/bin/dbus-monitor
/apps/mousepad/local/bin/dbus-send
/apps/mousepad/local/bin/dbus-uuidgen
/apps/mousepad/local/bin/xfconf-query
/apps/mousepad/local/etc/bash_completion.d/dbus-bash-completion.sh
/apps/mousepad/local/etc/init.d/dbus
/apps/mousepad/local/lib/libEGL.so
/apps/mousepad/local/lib/libEGL.so.1
/apps/mousepad/local/lib/libEGL.so.1.0
/apps/mousepad/local/lib/libFS.so
/apps/mousepad/local/lib/libFS.so.6
/apps/mousepad/local/lib/libFS.so.6.0.0
/apps/mousepad/local/lib/libXRes.so
/apps/mousepad/local/lib/libXRes.so.1
/apps/mousepad/local/lib/libXRes.so.1.0.0
/apps/mousepad/local/lib/libXaw.so
/apps/mousepad/local/lib/libXaw.so.6
/apps/mousepad/local/lib/libXaw.so.7
/apps/mousepad/local/lib/libXaw6.so
/apps/mousepad/local/lib/libXaw6.so.6
/apps/mousepad/local/lib/libXaw6.so.6.0.1
/apps/mousepad/local/lib/libXaw7.so
/apps/mousepad/local/lib/libXaw7.so.7
/apps/mousepad/local/lib/libXaw7.so.7.0.0
/apps/mousepad/local/lib/libXcomposite.so
/apps/mousepad/local/lib/libXcomposite.so.1
/apps/mousepad/local/lib/libXcomposite.so.1.0.0
/apps/mousepad/local/lib/libXcursor.so
/apps/mousepad/local/lib/libXcursor.so.1
/apps/mousepad/local/lib/libXcursor.so.1.0.2
/apps/mousepad/local/lib/libXdamage.so
/apps/mousepad/local/lib/libXdamage.so.1
/apps/mousepad/local/lib/libXdamage.so.1.1.0
/apps/mousepad/local/lib/libXfixes.so
/apps/mousepad/local/lib/libXfixes.so.3
/apps/mousepad/local/lib/libXfixes.so.3.1.0
/apps/mousepad/local/lib/libXinerama.so
/apps/mousepad/local/lib/libXinerama.so.1
/apps/mousepad/local/lib/libXinerama.so.1.0.0
/apps/mousepad/local/lib/libXss.so
/apps/mousepad/local/lib/libXss.so.1
/apps/mousepad/local/lib/libXss.so.1.0.0
/apps/mousepad/local/lib/libXtst.so
/apps/mousepad/local/lib/libXtst.so.6
/apps/mousepad/local/lib/libXtst.so.6.1.0
/apps/mousepad/local/lib/libXv.so
/apps/mousepad/local/lib/libXv.so.1
/apps/mousepad/local/lib/libXv.so.1.0.0
/apps/mousepad/local/lib/libXvMC.so
/apps/mousepad/local/lib/libXvMC.so.1
/apps/mousepad/local/lib/libXvMC.so.1.0.0
/apps/mousepad/local/lib/libXvMCW.so
/apps/mousepad/local/lib/libXvMCW.so.1
/apps/mousepad/local/lib/libXvMCW.so.1.0.0
/apps/mousepad/local/lib/libXxf86dga.so
/apps/mousepad/local/lib/libXxf86dga.so.1
/apps/mousepad/local/lib/libXxf86dga.so.1.0.0
/apps/mousepad/local/lib/libXxf86vm.so
/apps/mousepad/local/lib/libXxf86vm.so.1
/apps/mousepad/local/lib/libXxf86vm.so.1.0.0
/apps/mousepad/local/lib/libdbus-1.so
/apps/mousepad/local/lib/libdbus-1.so.3
/apps/mousepad/local/lib/libdbus-1.so.3.7.2
/apps/mousepad/local/lib/libdbus-glib-1.so
/apps/mousepad/local/lib/libdbus-glib-1.so.2
/apps/mousepad/local/lib/libdbus-glib-1.so.2.2.2
/apps/mousepad/local/lib/libfontenc.so
/apps/mousepad/local/lib/libfontenc.so.1
/apps/mousepad/local/lib/libfontenc.so.1.0.0
/apps/mousepad/local/lib/libglade/2.0/libxfce4.a
/apps/mousepad/local/lib/libglade/2.0/libxfce4.la
/apps/mousepad/local/lib/libglade/2.0/libxfce4.so
/apps/mousepad/local/lib/libglade-2.0.so
/apps/mousepad/local/lib/libglade-2.0.so.0
/apps/mousepad/local/lib/libglade-2.0.so.0.0.7
/apps/mousepad/local/lib/libstartup-notification-1.so
/apps/mousepad/local/lib/libstartup-notification-1.so.0
/apps/mousepad/local/lib/libstartup-notification-1.so.0.0.0
/apps/mousepad/local/lib/libxcb-atom.so
/apps/mousepad/local/lib/libxcb-atom.so.1
/apps/mousepad/local/lib/libxcb-atom.so.1.0.0
/apps/mousepad/local/lib/libxcb-aux.so
/apps/mousepad/local/lib/libxcb-aux.so.0
/apps/mousepad/local/lib/libxcb-aux.so.0.0.0
/apps/mousepad/local/lib/libxcb-composite.so
/apps/mousepad/local/lib/libxcb-composite.so.0
/apps/mousepad/local/lib/libxcb-composite.so.0.0.0
/apps/mousepad/local/lib/libxcb-damage.so
/apps/mousepad/local/lib/libxcb-damage.so.0
/apps/mousepad/local/lib/libxcb-damage.so.0.0.0
/apps/mousepad/local/lib/libxcb-dpms.so
/apps/mousepad/local/lib/libxcb-dpms.so.0
/apps/mousepad/local/lib/libxcb-dpms.so.0.0.0
/apps/mousepad/local/lib/libxcb-dri2.so
/apps/mousepad/local/lib/libxcb-dri2.so.0
/apps/mousepad/local/lib/libxcb-dri2.so.0.0.0
/apps/mousepad/local/lib/libxcb-event.so
/apps/mousepad/local/lib/libxcb-event.so.1
/apps/mousepad/local/lib/libxcb-event.so.1.0.0
/apps/mousepad/local/lib/libxcb-glx.so
/apps/mousepad/local/lib/libxcb-glx.so.0
/apps/mousepad/local/lib/libxcb-glx.so.0.0.0
/apps/mousepad/local/lib/libxcb-icccm.so
/apps/mousepad/local/lib/libxcb-icccm.so.1
/apps/mousepad/local/lib/libxcb-icccm.so.1.0.0
/apps/mousepad/local/lib/libxcb-image.so
/apps/mousepad/local/lib/libxcb-image.so.0
/apps/mousepad/local/lib/libxcb-image.so.0.0.0
/apps/mousepad/local/lib/libxcb-keysyms.so
/apps/mousepad/local/lib/libxcb-keysyms.so.1
/apps/mousepad/local/lib/libxcb-keysyms.so.1.0.0
/apps/mousepad/local/lib/libxcb-property.so
/apps/mousepad/local/lib/libxcb-property.so.1
/apps/mousepad/local/lib/libxcb-property.so.1.0.0
/apps/mousepad/local/lib/libxcb-randr.so
/apps/mousepad/local/lib/libxcb-randr.so.0
/apps/mousepad/local/lib/libxcb-randr.so.0.1.0
/apps/mousepad/local/lib/libxcb-record.so
/apps/mousepad/local/lib/libxcb-record.so.0
/apps/mousepad/local/lib/libxcb-record.so.0.0.0
/apps/mousepad/local/lib/libxcb-render-util.so
/apps/mousepad/local/lib/libxcb-render-util.so.0
/apps/mousepad/local/lib/libxcb-render-util.so.0.0.0
/apps/mousepad/local/lib/libxcb-render.so
/apps/mousepad/local/lib/libxcb-render.so.0
/apps/mousepad/local/lib/libxcb-render.so.0.0.0
/apps/mousepad/local/lib/libxcb-reply.so
/apps/mousepad/local/lib/libxcb-reply.so.1
/apps/mousepad/local/lib/libxcb-reply.so.1.0.0
/apps/mousepad/local/lib/libxcb-res.so
/apps/mousepad/local/lib/libxcb-res.so.0
/apps/mousepad/local/lib/libxcb-res.so.0.0.0
/apps/mousepad/local/lib/libxcb-screensaver.so
/apps/mousepad/local/lib/libxcb-screensaver.so.0
/apps/mousepad/local/lib/libxcb-screensaver.so.0.0.0
/apps/mousepad/local/lib/libxcb-shape.so
/apps/mousepad/local/lib/libxcb-shape.so.0
/apps/mousepad/local/lib/libxcb-shape.so.0.0.0
/apps/mousepad/local/lib/libxcb-shm.so
/apps/mousepad/local/lib/libxcb-shm.so.0
/apps/mousepad/local/lib/libxcb-shm.so.0.0.0
/apps/mousepad/local/lib/libxcb-sync.so
/apps/mousepad/local/lib/libxcb-sync.so.0
/apps/mousepad/local/lib/libxcb-sync.so.0.0.0
/apps/mousepad/local/lib/libxcb-xevie.so
/apps/mousepad/local/lib/libxcb-xevie.so.0
/apps/mousepad/local/lib/libxcb-xevie.so.0.0.0
/apps/mousepad/local/lib/libxcb-xf86dri.so
/apps/mousepad/local/lib/libxcb-xf86dri.so.0
/apps/mousepad/local/lib/libxcb-xf86dri.so.0.0.0
/apps/mousepad/local/lib/libxcb-xfixes.so
/apps/mousepad/local/lib/libxcb-xfixes.so.0
/apps/mousepad/local/lib/libxcb-xfixes.so.0.0.0
/apps/mousepad/local/lib/libxcb-xinerama.so
/apps/mousepad/local/lib/libxcb-xinerama.so.0
/apps/mousepad/local/lib/libxcb-xinerama.so.0.0.0
/apps/mousepad/local/lib/libxcb-xprint.so
/apps/mousepad/local/lib/libxcb-xprint.so.0
/apps/mousepad/local/lib/libxcb-xprint.so.0.0.0
/apps/mousepad/local/lib/libxcb-xtest.so
/apps/mousepad/local/lib/libxcb-xtest.so.0
/apps/mousepad/local/lib/libxcb-xtest.so.0.0.0
/apps/mousepad/local/lib/libxcb-xv.so
/apps/mousepad/local/lib/libxcb-xv.so.0
/apps/mousepad/local/lib/libxcb-xv.so.0.0.0
/apps/mousepad/local/lib/libxcb-xvmc.so
/apps/mousepad/local/lib/libxcb-xvmc.so.0
/apps/mousepad/local/lib/libxcb-xvmc.so.0.0.0
/apps/mousepad/local/lib/libxcb.so
/apps/mousepad/local/lib/libxcb.so.1
/apps/mousepad/local/lib/libxcb.so.1.1.0
/apps/mousepad/local/lib/libxfce4util.so
/apps/mousepad/local/lib/libxfce4util.so.4
/apps/mousepad/local/lib/libxfce4util.so.4.1.1
/apps/mousepad/local/lib/libxfce4util.so.6
/apps/mousepad/local/lib/libxfce4util.so.6.0.0
/apps/mousepad/local/lib/libxfcegui4.so
/apps/mousepad/local/lib/libxfcegui4.so.4
/apps/mousepad/local/lib/libxfcegui4.so.4.3.0
/apps/mousepad/local/lib/libxfconf-0.so
/apps/mousepad/local/lib/libxfconf-0.so.2
/apps/mousepad/local/lib/libxfconf-0.so.2.0.0
/apps/mousepad/local/lib/libxkbfile.so
/apps/mousepad/local/lib/libxkbfile.so.1
/apps/mousepad/local/lib/libxkbfile.so.1.0.2
/apps/mousepad/local/lib/xfce4/xfconf/xfconfd
/apps/mousepad/local/libexec/dbus-bash-completion-helper
/apps/mousepad/local/libexec/dbus-daemon-launch-helper
/apps/mousepad/local/sbin/xfce4-kiosk-query
/apps/mousepad/local/share/dbus/files/session.conf
/apps/mousepad/local/share/dbus/files/system.conf
/apps/mousepad/local/share/dbus-1/services/org.xfce.Xfconf.service
/apps/mousepad/local/share/icons/hicolor/48x48/apps/xfce-edit.png
/apps/mousepad/local/share/icons/hicolor/48x48/apps/xfce-filemanager.png
/apps/mousepad/local/share/icons/hicolor/48x48/apps/xfce-games.png
/apps/mousepad/local/share/icons/hicolor/48x48/apps/xfce-graphics.png
/apps/mousepad/local/share/icons/hicolor/48x48/apps/xfce-internet.png
/apps/mousepad/local/share/icons/hicolor/48x48/apps/xfce-man.png
/apps/mousepad/local/share/icons/hicolor/48x48/apps/xfce-multimedia.png
/apps/mousepad/local/share/icons/hicolor/48x48/apps/xfce-printer.png
/apps/mousepad/local/share/icons/hicolor/48x48/apps/xfce-schedule.png
/apps/mousepad/local/share/icons/hicolor/48x48/apps/xfce-sound.png
/apps/mousepad/local/share/icons/hicolor/48x48/apps/xfce-terminal.png
/apps/mousepad/local/share/icons/hicolor/48x48/apps/xfce-unknown.png
/apps/mousepad/local/share/icons/hicolor/48x48/apps/xfce-utils.png
/apps/mousepad/local/share/icons/hicolor/scalable/apps/xfce-filemanager.svg
/apps/mousepad/local/share/locale/am/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/ar/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/be/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/ca/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/cs/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/cy/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/da/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/de/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/dz/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/el/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/en_GB/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/eo/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/es/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/et/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/eu/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/fi/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/fr/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/gl/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/gu/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/he/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/hu/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/id/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/it/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/ja/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/ka/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/ko/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/ku/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/lt/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/lv/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/mk/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/nb/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/nl/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/pa/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/pl/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/pt_BR/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/pt_PT/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/ru/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/si/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/sk/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/sl/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/sq/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/sv/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/ta/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/tr/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/uk/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/ur/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/zh_CN/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/locale/zh_TW/LC_MESSAGES/mousepad.mo
/apps/mousepad/local/share/xml/libglade/glade-2.0.dtd
/apps/mousepad/local/tce.installed/Xlibs
/apps/mousepad/local/tce.installed/dbus
/apps/mousepad/mousepad
/apps/mousepad/share/applications/mousepad.desktop
/apps/mousepad/share/pixmaps/mousepad.png
